# roles/setup/files/tarsnap/tarsnap.timer
[Unit]
Description=Tarsnap Backups Timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
