# roles/setup/files/password-store/password-store-sync.service
[Service]
CPUSchedulingPolicy=idle
ExecStart=/usr/local/bin/password-store-sync.sh
IOSchedulingClass=idle

[Unit]
Description=Password store sync
